# scripts/nexus-server.service
# nexus-server.service — Copy vào /etc/systemd/system/
# sudo cp scripts/nexus-server.service /etc/systemd/system/
# sudo systemctl daemon-reload && sudo systemctl enable nexus-server

[Unit]
Description=Nexus Isekai Game Server
After=network.target mysql.service
Requires=mysql.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/nexus-isekai/server

ExecStart=/usr/bin/java \
    -Xmx2g -Xms512m \
    -XX:+UseG1GC \
    -XX:MaxGCPauseMillis=50 \
    -XX:+UseStringDeduplication \
    -Djava.awt.headless=true \
    -jar target/nexus-isekai-server-1.0.jar

Restart=on-failure
RestartSec=15
StandardOutput=append:/var/log/nexus/server.log
StandardError=append:/var/log/nexus/server-error.log
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
